# codio/base/init/codio.conf
description "Codio startup script"

start on startup

limit nproc 512 512
limit nofile 512 512

post-start script
    exec /usr/sbin/codio-start
end script
