# base/init/codio.conf
description "Codio startup script"

start on startup

post-start script
    exec /usr/sbin/codio-start
end script
